# 5-multistuff/uber-cool-app/deploy/uber-cool-app.conf
[program:uber-cool-app]
command=/usr/bin/node /app/app.js

